# SOURCES/init.d-graylog2-web-interface
#!/bin/bash

#
# graylog2-web-interface:   Web interface for the graylog2 message collector
#
# chkconfig: - 98 02
# description:  This daemon provides a front-end web interface for the Graylog2 
#               message collector
#

CMD=$1

# resolve links - $0 may be a softlink
GRAYLOG2CTL="$0"

while [ -h "$GRAYLOG2CTL" ]; do
    ls=`ls -ld "$GRAYLOG2CTL"`
    link=`expr "$ls" : '.*-> \(.*\)$'`
    if expr "$link" : '/.*' > /dev/null; then
        GRAYLOG2CTL="$link"
    else
        GRAYLOG2CTL=`dirname "$GRAYLOG2CTL"`/"$link"
    fi
done

#####################
## Any changes to settings make in /etc/sysconfig/graylog2-web-interface
#####################
JAVA_CMD=`which java 2> /dev/null`

GRAYLOG2CTL_DIR=/opt/graylog2-web-interface
GRAYLOG2_PID=/opt/graylog2-web-interface/RUNNING_PID
GRAYLOG2_USER=root
LOG_FILE=/var/log/graylog2/graylog2-web-interface.log
ADDITIONAL_JAVA_OPTS="-XX:ReservedCodeCacheSize=128m"

lib_dir="${GRAYLOG2CTL_DIR}/lib"
CLASSPATH="${lib_dir}/*"

MIN_MEM=1024m
MAX_MEM=1024m
MAX_PERM_SIZE=256m
#####################
## End changes
#####################
[ -f /etc/sysconfig/graylog2-web-interface ] && . /etc/sysconfig/graylog2-web-interface
[ -f /etc/rc.d/init.d/functions ] && . /etc/rc.d/init.d/functions

RETVAL=0

# Make sure that the java exec is found
 
if [[ "" =~ "$JAVA_CMD" || ! -x $JAVA_CMD ]] ; then
    echo "ERROR: Java not found! Exiting."
    exit
fi

start() {
    echo -n "Starting graylog2-web-interface ..."
    cd "$GRAYLOG2CTL_DIR/"
    su -s /bin/sh $GRAYLOG2_USER -c "${JAVA_CMD} -Xms${MIN_MEM} -Xmx${MAX_MEM} \
                                     -XX:MaxPermSize=${MAX_PERM_SIZE} ${ADDITIONAL_JAVA_OPTS} \
                                     -Duser.dir=/opt/graylog2-web-interface \
                                     -Dpidfile.path=${GRAYLOG2_PID} \
                                     -cp \"$CLASSPATH\" play.core.server.NettyServer >> ${LOG_FILE} 2>&1 &"
    RETVAL=$?

    if [ $RETVAL -eq 0 ]
    then
        echo_success
    else
        echo_failure
    fi

    echo
}

stop() {
    PID=$(get_pid)
    echo -n "Stopping graylog2-web-interface ($PID) ..."

    if kill $PID; then
        rm ${GRAYLOG2_PID}
        echo_success
    else
        echo_failure
    fi

    echo
}

restart() {
    echo "Restarting graylog2-web-interface ..."
    stop
    start
}

status() {
    pid=$(get_pid)
    if [ ! -z $pid ]; then
        if pid_running $pid; then
            echo "graylog2-web-interface running as pid $pid"
            return 0
        else
            echo "Stale pid file with $pid - removing..."
            rm ${GRAYLOG2_PID}
        fi
    fi
 
    echo "graylog2-web-interface not running"
}
get_pid() {
    cat ${GRAYLOG2_PID} 2> /dev/null
}
 
pid_running() {
    kill -0 $1 2> /dev/null
}

case "$CMD" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        restart
        ;;
    status)
        status
        ;;
    *)
        echo "Usage $0 {start|stop|restart|status}"
        RETVAL=1
esac
